# getpix.timer
[Unit]
Description= check every minute

[Timer]
Persistent=true
OnBootSec=1min
OnCalendar=*:0,15,30,45

[Install]
WantedBy=timers.target
