# getpix.timer
[Unit]
Description= check every minute

[Timer]
Persistent=true
OnBootSec=1min
OnCalendar=minutely

[Install]
WantedBy=timers.target
